# src/filesystem/etc/rustysd/services/default.target
[Unit]
Description=Things are good when this is running